# templates/cron.d/adcli_update
# $FreeBSD$
#
SHELL=/bin/sh
PATH=/etc:/bin:/sbin:/usr/bin:/usr/sbin

#minute hour    mday    month   wday    who     command
1	1	*	*	*	root	/bin/sh /var/db/samba4/adcli_update.sh
