# templates/periodic/weekly/600.machine_password_reset
#!/bin/sh -
#
# $FreeBSD: head/audio/ampache/files/500.ampache-cron.in 564645 2021-02-07 20:36:26Z crees $
#

if [ -r /etc/defaults/periodic.conf ]; then
	. /etc/defaults/periodic.conf
	source_periodic_confs
fi

rc=0

case "${weekly_machine_password_reset_enable:-no}" in
[Nn][Oo])
	;;
*)
	/usr/local/bin/net ads changetrustpw
	;;
esac

exit $rc
